# configs/ardentlinux-doctor/airootfs/etc/systemd/system/getty@tty2.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear %i 38400 $TERM
