# configs/ardentlinux-amiya/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin doctor --noclear %i 38400 linux
