# configs/ardentlinux-doctor/airootfs/etc/systemd/system/getty@tty3.service.d/autologin.conf
[Service]
Environment=XDG_SESSION_TYPE=x11
ExecStart=
ExecStart=-/sbin/agetty --autologin doctor --noclear %i 38400 $TERM
